# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── 1. Chromium headless com CDP (sobe primeiro) ─────────────────────
[program:chromium]
command=/app/scripts/start-browser.sh
user=1001
autostart=true
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/data/logs/chromium.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/data/logs/chromium-err.log
stderr_logfile_maxbytes=5MB
priority=10

# ── 2. Nginx proxy CDP (9223 → 9222) ─────────────────────────────────
# Proxy para o Chromium CDP, permite que o openclaw acesse
# o browser via BROWSER_CDP_URL=http://127.0.0.1:9223
[program:nginx-cdp]
command=nginx -c /etc/nginx/cdp-proxy.conf -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/data/logs/nginx-cdp.log
stdout_logfile_maxbytes=2MB
stderr_logfile=/data/logs/nginx-cdp-err.log
stderr_logfile_maxbytes=2MB
priority=20

# ── 3. Nginx principal (reverse proxy :PORT → openclaw gateway) ───────
[program:nginx-main]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=2
stdout_logfile=/data/logs/nginx-main.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/data/logs/nginx-main-err.log
stderr_logfile_maxbytes=2MB
priority=30

# ── 4. OpenClaw Gateway ───────────────────────────────────────────────
# O binary openclaw está no PATH (instalado pelo coollabsio/openclaw base)
# configure.js já rodou no entrypoint — apenas inicia o gateway
[program:openclaw]
command=openclaw gateway
environment=HOME="/root",OPENCLAW_HOME="/data",OPENCLAW_STATE_DIR="/data/.openclaw",OPENCLAW_CONFIG_PATH="/data/.openclaw/openclaw.json"
autostart=true
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/data/logs/openclaw.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/openclaw-err.log
stderr_logfile_maxbytes=10MB
priority=40
